# deploy/systemd/discord-gateway.service
[Unit]
Description=Discord gateway (realtime two-way conversation / ears) for embodied-claude
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ytsuhako/projects/embodied-claude
ExecStart=/home/ytsuhako/.local/bin/uv --directory /home/ytsuhako/projects/embodied-claude/discord-mcp run discord-gateway
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
